# openwrt/luci-app-kwrtmgrd/root/etc/uci-defaults/40_luci-kwrtmgrd
#!/bin/sh
# 刷新 LuCI 菜单/模块缓存并重载 rpcd，让 kwrtmgr 菜单与 ACL 立即生效
rm -f /tmp/luci-indexcache* 2>/dev/null
rm -rf /tmp/luci-modulecache 2>/dev/null
/etc/init.d/rpcd reload 2>/dev/null
exit 0
