# openwrt/files/etc/init.d/kwrtmgrd
#!/bin/sh /etc/rc.common
# =============================================================================
# kwrtmgrd — OpenWrt procd 服务脚本
#
#   由 ipk 包安装到 /etc/init.d/kwrtmgrd，受 procd 守护。
#   配置读自 UCI: /etc/config/kwrtmgrd (main 节)，转成 KWRTNET_* 环境变量注入。
#   守护进程本身只认环境变量（见 internal/appcfg），故二进制零改动。
#
#   子命令 start/stop/restart/reload/enable/disable/status 由 /etc/rc.common 提供：
#     /etc/init.d/kwrtmgrd start      启动
#     /etc/init.d/kwrtmgrd enable      开机自启
#   修改配置后:  uci commit kwrtmgrd; /etc/init.d/kwrtmgrd restart
# =============================================================================

USE_PROCD=1
# 网络起来后再启动；停止时尽早
START=95
STOP=01

PROG=/usr/bin/kwrtmgrd
CONF=kwrtmgrd

# ---------------------------------------------------------------------------
# 生成强随机 16 进制令牌 (48 字符)。取不到随机源时返回空串 (由 ensure_token 拒绝启动)。
#
# 注意：busybox 的 applet 集合按固件裁剪，不能假定某个工具一定存在
# (实测 ImmortalWrt 24.10 x86/64 默认无 od applet -> 旧实现 token 恒为空 -> 拒绝启动)。
# 故按可用性逐级回退，任一产出 >=32 字符即采用：
#   ① od        ② hexdump   ③ sha256sum  ④ md5sum(两段拼接)
#   ⑤ 内核 uuid (/proc/sys/kernel/random/uuid，几乎必有)
#   ⑥ 直接从 /dev/urandom 流筛十六进制字符 (tr/head 为核心 applet)
# 全部失败 (无随机源) 才返回空串，由 ensure_token 拒绝以弱令牌启动。
# ---------------------------------------------------------------------------
gen_token() {
	local t=''
	if [ -r /dev/urandom ] && command -v od >/dev/null 2>&1; then
		t="$(dd if=/dev/urandom bs=1 count=24 2>/dev/null | od -An -tx1 2>/dev/null | tr -d ' \n')"
	fi
	if [ "${#t}" -lt 32 ] && [ -r /dev/urandom ] && command -v hexdump >/dev/null 2>&1; then
		t="$(dd if=/dev/urandom bs=1 count=24 2>/dev/null | hexdump -v -e '/1 "%02x"' 2>/dev/null)"
	fi
	if [ "${#t}" -lt 32 ] && [ -r /dev/urandom ] && command -v sha256sum >/dev/null 2>&1; then
		t="$(dd if=/dev/urandom bs=1 count=64 2>/dev/null | sha256sum 2>/dev/null | cut -c1-48)"
	fi
	if [ "${#t}" -lt 32 ] && [ -r /dev/urandom ] && command -v md5sum >/dev/null 2>&1; then
		t="$(dd if=/dev/urandom bs=1 count=64 2>/dev/null | md5sum 2>/dev/null | cut -c1-32)$(dd if=/dev/urandom bs=1 count=64 2>/dev/null | md5sum 2>/dev/null | cut -c1-16)"
	fi
	if [ "${#t}" -lt 32 ] && [ -r /proc/sys/kernel/random/uuid ]; then
		t="$(cat /proc/sys/kernel/random/uuid 2>/dev/null | tr -d '-')$(cat /proc/sys/kernel/random/uuid 2>/dev/null | tr -d '-')"
		t="$(printf '%s' "$t" | cut -c1-48)"
	fi
	if [ "${#t}" -lt 32 ] && [ -r /dev/urandom ]; then
		t="$(tr -dc 'a-f0-9' < /dev/urandom 2>/dev/null | head -c 48 2>/dev/null)"
	fi
	printf '%s' "$t"
}

# ---------------------------------------------------------------------------
# 确保存在登录令牌。守护进程要求 KWRTNET_API_TOKEN 非空，否则启动即退出，
# 因此首次启动 (UCI token 为空) 时自动生成一个并落盘。幂等。
# 安全优先：拿不到足够随机度时宁可拒绝启动 (返回非 0)，绝不退化成可预测弱令牌。
# ---------------------------------------------------------------------------
ensure_token() {
	local t
	config_load "$CONF"
	config_get t main token ''
	[ -n "$t" ] && return 0

	t="$(gen_token)"
	if [ "${#t}" -lt 32 ]; then
		logger -t kwrtmgrd "ERROR: 无法生成安全随机令牌 (系统缺少 od/hexdump/sha256sum/md5sum 等十六进制工具，且无 /dev/urandom 与内核 uuid 随机源)，拒绝以弱令牌启动。请手动设置: uci set kwrtmgrd.main.token=<你的令牌>; uci commit kwrtmgrd; /etc/init.d/kwrtmgrd restart"
		return 1
	fi
	uci -q set "${CONF}.main.token=$t"
	uci -q commit "$CONF"
	logger -t kwrtmgrd "首次启动：已生成随机 API 令牌 (uci get kwrtmgrd.main.token 查看)"
}

start_service() {
	local enabled http_addr token data_dir log_level docs_enabled cors_origins self_update

	config_load "$CONF"
	config_get_bool enabled main enabled 1
	if [ "$enabled" != "1" ]; then
		# enabled 是「持久运行状态」：网页点停止会置 0。无论系统重启还是升级，
		# 只要它是 0 就保持停止（除非勾选了开机强制自启，由网页把它重新置 1）。
		logger -t kwrtmgrd "当前为停止状态 (option enabled '0')，不启动；如需运行请在网页点「启动」或勾选「开机自动启动」"
		return 0
	fi

	# 先确保有令牌 (失败则拒绝启动)，再重新读取拿到最新值
	ensure_token || return 1

	config_load "$CONF"
	config_get      http_addr    main http_addr    ':18080'
	config_get      token        main token        ''
	config_get      data_dir     main data_dir     '/usr/lib/kwrtmgrd'
	config_get      log_level    main log_level    'info'
	config_get_bool docs_enabled main docs_enabled 1
	config_get      cors_origins main cors_origins '*'
	# OpenWrt 上 Web 端一键自更新会覆盖 opkg/apk 管理的二进制，与包管理器冲突，
	# 默认禁用；升级请走 opkg/apk。可在 UCI 显式打开 (不推荐)。
	config_get_bool self_update  main self_update  0

	# 令牌兜底校验：为空则拒绝启动 (否则守护进程读不到 token 会自行退出)
	[ -n "$token" ] || { logger -t kwrtmgrd "ERROR: API 令牌为空，拒绝启动"; return 1; }

	[ "$docs_enabled" = "1" ] && docs_enabled="true" || docs_enabled="false"
	[ "$self_update" = "1" ]  && self_update="true"  || self_update="false"

	# 数据目录需持久 (OpenWrt 上 /var -> /tmp 是 tmpfs，重启即丢，切勿用)；
	# 不可写时仅告警不阻断，提示改到可写持久路径 (或先配 extroot)
	mkdir -p "$data_dir" 2>/dev/null
	if touch "$data_dir/.wtest" 2>/dev/null; then
		rm -f "$data_dir/.wtest" 2>/dev/null
	else
		logger -t kwrtmgrd "WARNING: 数据目录不可写: $data_dir，请 uci set kwrtmgrd.main.data_dir=<可写持久路径> 后重启"
	fi

	# 二进制由 all 包安装时经 kwrtmgrd-fetch 联网拉取；若缺失（如装时无网）则拒绝启动并提示
	[ -x "$PROG" ] || {
		logger -t kwrtmgrd "ERROR: 二进制未安装 ($PROG)，请联网后执行: kwrtmgrd-fetch  然后 /etc/init.d/kwrtmgrd start"
		return 1
	}

	procd_open_instance
	procd_set_param command "$PROG" serve
	procd_set_param env \
		KWRTNET_API_TOKEN="$token" \
		KWRTNET_HTTP_ADDR="$http_addr" \
		KWRTNET_DATA_DIR="$data_dir" \
		KWRTNET_LOG_LEVEL="$log_level" \
		KWRTNET_DOCS_ENABLED="$docs_enabled" \
		KWRTNET_CORS_ORIGINS="$cors_origins" \
		KWRTNET_SELF_UPDATE_ENABLED="$self_update"
	# 崩溃自动拉起：3600s 内异常退出算崩溃，重试 5 次；之后退出则无限重启
	procd_set_param respawn "${respawn_threshold:-3600}" "${respawn_timeout:-5}" "${respawn_retry:-5}"
	# 把 stdout/stderr 接入 logd，可用 logread 查看
	procd_set_param stdout 1
	procd_set_param stderr 1
	# 配置文件变更触发实例重载
	procd_set_param file /etc/config/kwrtmgrd
	procd_close_instance
}

service_triggers() {
	# 监听 config.change 事件，配置变更后自动 reload
	procd_add_reload_trigger "kwrtmgrd"
}
